# pi-frame_update.timer
[Unit]
Description=Updating the pi-frame on boot.

[Timer]
OnBootSec=90
OnUnitActiveSec=5m

[Install]
WantedBy=timers.target
